# Archive/with2PortVMs/Ansible/roles/ndnping/files/etc_init_ndnping.conf
description "NDN TLV ping server"
author      "Alexander Afanasyev <alexander.afanasyev@ucla.edu>"

start on started  nfd
stop  on stopping nfd

respawn
respawn limit unlimited

setuid ndn
setgid ndn

script
    FLAGS=""
    PREFIX="/localhost/ndnpingserver"
    if [ -f /usr/local/etc/default/ndnping ]; then
        . /usr/local/etc/default/ndnping
    fi
    export HOME=/usr/local/var/lib/ndn/ndnping
    exec /usr/local/bin/ndnpingserver "$PREFIX" $FLAGS
end script
